# ops/systemd/sentinel-sensor@.service
[Unit]
Description=Sentinel NetLab Sensor (iface %i)
After=network-online.target syslog.target
Wants=network-online.target
# Ensure config exists before starting
ConditionPathExists=/etc/sentinel/config.yaml

[Service]
Type=simple

# Run as unprivileged user 'sentinel' (created in setup). Systemd will grant needed capabilities.
User=sentinel
Group=sentinel

# Grant necessary network capabilities for monitor mode / raw sockets.
# AmbientCapabilities allows the service process to have CAP_NET_RAW and CAP_NET_ADMIN.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Avoid escalation, but keep capabilities for network ops
NoNewPrivileges=true

# Environment file for secrets (BEARER_TOKEN, etc.)
EnvironmentFile=/etc/sentinel/env

# Pre-start: optional helper to put iface into monitor mode (script must be made idempotent & safe).
# It receives the interface as %i (instance name)
ExecStartPre=/usr/local/bin/sentinel-ensure-monitor-mode %i

# Main command: adjust to your venv/python path and CLI flags
ExecStart=/opt/sentinel/venv/bin/python /opt/sentinel/sensor/cli.py --config /etc/sentinel/config.yaml --iface %i

# Graceful shutdown
KillMode=control-group
TimeoutStopSec=20

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Resource limits
LimitNOFILE=4096
LimitNPROC=512
LimitCORE=infinity

# Minimal filesystem protections
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Logging (keep journald default)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-sensor-%i

[Install]
WantedBy=multi-user.target
