# ops/systemd/wifi-scanner.service
[Unit]
Description=Sentinel NetLab WiFi Sensor API
Documentation=https://github.com/Anduong1200/sentinel-netlab
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sentinel-netlab/sensor

# Environment
Environment="PATH=/opt/sentinel-netlab/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="WIFI_SCANNER_INTERFACE=wlan0"
Environment="WIFI_SCANNER_API_KEY=sentinel-dev-2024"
Environment="CAPTURE_ENGINE=tshark"
Environment="ALLOW_ACTIVE_ATTACKS=false"

# Run as root for monitor mode access
# TODO: Switch to dedicated user after setting up capabilities
User=root
Group=root

# Application start
ExecStartPre=/bin/bash -c 'ip link set ${WIFI_SCANNER_INTERFACE} down 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'iw dev ${WIFI_SCANNER_INTERFACE} set type monitor 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'ip link set ${WIFI_SCANNER_INTERFACE} up 2>/dev/null || true'
ExecStart=/opt/sentinel-netlab/venv/bin/python api_server.py

# Restart policy
Restart=always
RestartSec=5
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-scanner

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
